# data/fleet-commander-controller.service.in
[Unit]
Description=Fleet Commander virtual session controller
After=network.target remote-fs.target nss-lookup.target

[Service]
ExecStart=@PRIVLIBEXECDIR@/fleet-commander-controller.py

StandardOutput=syslog
StandardError=inherit
