# data/fleet-commander-controller.service.in
[Unit]
Description=Fleet Commander virtual session controller
After=network.target remote-fs.target nss-lookup.target

[Service]
Environment=PYTHONPATH=@FCPYTHONDIR@
ExecStart=@PYTHON@ -m fleetcommander.controller --configuration @XDGCONFIGDIR@/fleet-commander-controller.conf

StandardOutput=syslog
StandardError=inherit

[Install]
Alias=fleet-commander-controller.service
