# data/fleet-commander-vnc-session.service.in
[Unit]
Description=Fleet Commander VNC Session
After=systemd-user-sessions.service getty@tty1.service plymouth-quit.service

[Service]
ExecStart=@RUNUSER@ -l fleet-commander -c "@PRIVLIBEXECDIR@/fleet-commander-xvnc.sh"
ExecStartPost=@RUNUSER@ -l fleet-commander -c "@PRIVLIBEXECDIR@/fleet-commander-xinitrc.sh"

ExecStopPost=-@PKILL@ -u fleet-commander
ExecStopPost=-@PKILL@ -9 -u fleet-commander
ExecStopPost=-@RM@ -rf ~fleet-commander/.config/dconf/

#TODO: Use an overlayfs mount for the home directory

StandardOutput=syslog
StandardError=inherit
