# data/fleet-commander-admin.service.in
[Unit]
Description=Fleet Commander Administration Interface
After=network.target remote-fs.target nss-lookup.target

[Service]
User=fleet-commander-admin
ExecStart=@PRIVLIBEXECDIR@/fleet-commander-admin.py --configuration @XDGCONFIGDIR@/fleet-commander-admin.conf

StandardOutput=syslog
StandardError=inherit

[Install]
Alias=fleet-commander-admin.service
